# cert-monitor.service
[Unit]
Description=Certificate change monitor

[Service]
ExecStart=/usr/bin/cert-monitor
Restart=on-failure

[Install]
WantedBy=multi-user.target
